# cron_jobs/crontab
# Refesh host-ip pairs every minute
*/1 * * * * python3 /usr/local/bin/refresh_hosts.py
